# testdziej-loop.service
[Unit]
Description=Testdziej Question Generation Loop
After=network.target

[Service]
Type=simple
User=macryba
WorkingDirectory=/home/macryba/testdziej-questions
ExecStart=/usr/local/bin/claude loop --interval 30m --instructions .claude/instructions.md
Restart=always
RestartSec=60
StandardOutput=append:/home/macryba/testdziej-questions/logs/loop-service.log
StandardError=append:/home/macryba/testdziej-questions/logs/loop-service-error.log

[Install]
WantedBy=multi-user.target
